# data/systemd/mackes-remmina-sync.timer
[Unit]
Description=Mackes Remmina sync timer — runs every 5 minutes
Documentation=man:mackes(1)

[Timer]
# Fire 30 s after login (so the GUI is up first), then every 5 min.
OnBootSec=30s
OnUnitActiveSec=5min
# RandomizedDelaySec spreads the load if many peers wake at once
RandomizedDelaySec=20s
Persistent=false

[Install]
WantedBy=timers.target
