# data/systemd/mackesd-kdc-bridge.service
[Unit]
Description=Mackes Shell — KDE Connect mesh-mDNS bridge
Documentation=https://github.com/matthewmackes/MAP2-RELEASES
Documentation=file:///usr/share/mackes-shell/docs/PROJECT_WORKLIST.md#phase-13
# Phase 13.2 lock 2026-05-19: per-peer user-systemd unit, started
# by `mackes-panel.service`. Listens for `_kdeconnect._udp.local`
# via Avahi on the LAN; re-announces remote phones discovered
# through `~/QNM-Shared/<peer>/kdc/announce.jsonl` on the local
# LAN. Result: `kdeconnectd` sees every remote phone as if it
# were on the LAN.
After=graphical-session.target avahi-daemon.service
PartOf=graphical-session.target
Requires=avahi-daemon.service
ConditionEnvironment=DISPLAY

[Service]
Type=simple
ExecStart=/usr/bin/mackesd-kdc-bridge
Restart=on-failure
RestartSec=5
# The bridge holds no privileged state — runs as the session user.
Nice=10

[Install]
WantedBy=graphical-session.target
