# data/systemd/mackes-remmina-sync.service
[Unit]
Description=Mackes Remmina sync — auto-populate Remmina with mesh peers
Documentation=man:mackes(1)
After=tailscaled.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 -m mackes.remmina_sync --once
# Don't fail the unit on transient network issues — the timer will retry
SuccessExitStatus=0 1
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
