# data/systemd/mde-mesh-mount@.service
[Unit]
# MESHFS-4.1 (v5.0.0) — per-user FUSE mount of mesh-storage at the
# operator's XDG directories. Retargets GF-4.1 from GlusterFS to
# LizardFS. Instance `%i` is the XDG subdir name
# (Documents / Pictures / Music / Videos / Downloads).
#
# Mounts the LizardFS `mesh-storage` export's `%i` subfolder at
# `$HOME/%i` via `mfsmount`. The active master is always reached
# through the floating overlay VIP (10.42.0.1) — never a peer's
# underlay address — so the mount survives master failover
# transparently (mfsmount reconnects to the same VIP address).
#
# `~/Local/` is deliberately excluded from mesh-mounting per the
# MESHFS design: Local is the offline scratch area + private state
# that must never replicate. Only Documents/Pictures/Music/Videos/
# Downloads are mesh-backed.
#
# Not auto-enabled by RPM install — birthright's MESHFS XDG-mount
# step enables instances once uid:gid 1000:1000 (MESHFS-5.1) AND
# the LizardFS export is ready. Until that step ships, operators can
# manually enable individual instances after the fleet bootstraps:
#   systemctl --user enable --now mde-mesh-mount@Documents.service
Description=Mackes mesh-storage FUSE mount at %h/%i
Documentation=file:///usr/share/doc/mackes-shell/help/mesh-storage.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# mfsmount connects to the floating VIP (10.42.0.1, the active
# master address). mfssubfolder=/%i restricts the mount to this
# peer's Documents/Pictures/... slice of the shared export.
# The master connection is TCP/9419; mfsmount retries automatically
# so brief master-failover windows are transparent to the operator.
ExecStartPre=/usr/bin/mkdir -p %h/%i
ExecStart=/usr/bin/mfsmount \
    -o mfsmaster=10.42.0.1,mfssubfolder=/%i \
    %h/%i
ExecStop=/usr/bin/fusermount -u %h/%i

[Install]
WantedBy=default.target
