# data/systemd/greetd.service.d/role-gate.conf
# MDE role-gate drop-in for greetd.service (E1.5 / E1.3, 2026-06-05).
#
# The graphical greeter only starts on a Workstation-pinned box. On a
# Lighthouse/Server box the ExecCondition exits non-zero, so systemd SKIPS
# greetd (it does not FAIL it) and the machine boots to a usable headless
# console — never to a broken/blank greeter on a role that has no display
# (E1.5 acceptance #4; E1.3 acceptance #4 "role mismatch rejected at unit
# start, logged to the journal rather than starting a forbidden service").
#
# `mackesd role-gate --min-rank 2` (E1.3) exits 0 iff the pinned deployment
# rank is >= 2 (Workstation), else non-zero after journaling the conflict;
# it resolves the rank from /var/lib/mde/role.toml (unpinned -> Workstation,
# malformed -> Lighthouse, fail-closed). Belt-and-suspenders with the DM swap
# (which only enables greetd on Workstation in the first place): even an
# accidentally-enabled greetd on a non-Workstation refuses to start.
#
# Drop-in (MDE owns this file, NOT the Fedora greetd package's main unit), so
# there is no dual-ownership conflict at dnf-install time.

[Service]
ExecCondition=/usr/bin/mackesd role-gate --min-rank 2
